# .config/systemd/user/eset-nod32-mock.service
[Unit]
Description=esets_daemon mock

[Service]
ExecStart=%h/.local/bin/esets_daemon

[Install]
WantedBy=multi-user.target
